# pipeline/dogwatch-notify.service
[Unit]
Description=DogWatch MQTT event notifier
After=network-online.target
Wants=network-online.target

[Service]
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883
Environment=MQTT_TOPIC=dogwatch
# %h expands to the home directory of the user running this (user) unit —
# adjust the venv/script path below to match your install location.
ExecStart=%h/.openclaw/workspace/scripts/dogwatch-venv/bin/python %h/.openclaw/workspace/scripts/dogwatch-notify.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
